# systemd/simplebooth-kiosk.service
[Unit]
Description=SimpleBooth Kiosk
After=graphical.target
Wants=graphical.target

[Service]
Type=exec
User=admin
Group=admin
Environment=DISPLAY=:0
Environment=HOME=/home/admin
ExecStart=/home/admin/start_simplebooth.sh
Restart=always
RestartSec=3
TimeoutStartSec=60

[Install]
WantedBy=graphical.target
